# config/systemd/user/poweralertd.service.d/override.conf
#
[Service]
ExecStart=
ExecStart=/usr/bin/poweralertd -S
